# bin/armhf/etc/systemd/system/edge.service
[Unit]
Description=OpenEdge Middleware
 
[Service]
ExecStart=/home/debian/bin/armhf/edge -c /home/debian/bin/armhf/test1.config
WorkingDirectory=/home/debian/bin/armhf

 
[Install]
WantedBy=multi-user.target